# darwin/autoupgrade/com.lazybox.brew-autoupgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template plist — install.sh replaces @@HOME@@ with the actual home directory
  path before copying to ~/Library/LaunchAgents/.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lazybox.brew-autoupgrade</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@@HOME@@/.local/bin/brew-autoupgrade</string>
  </array>

  <!-- Run daily at 09:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Do NOT run immediately on load -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>@@HOME@@/Library/Logs/brew-autoupgrade.log</string>
  <key>StandardErrorPath</key>
  <string>@@HOME@@/Library/Logs/brew-autoupgrade.log</string>
</dict>
</plist>
